# roles/docker/templates/docker-prune.target.j2
[Unit]
Description=Docker prune target

[Install]
WantedBy=default.target
